# rabbit/rabbit.conf
[program:rabbitmq]
command=/bin/rabbit.sh
